# fc_temp/fc_wiz_postgres.service

[Unit]
Description = fc_wiz_postgres

[Service]
Type           = simple
User           = root
Group          = root
LimitNOFILE    = 4096
Restart        = always
RestartSec     = 10s
StandardOutput = append:/mnt/github/tjf/fc_logs/fc_wiz_postgres_access.log
StandardError  = append:/mnt/github/tjf/fc_logs/fc_wiz_postgres_error.log
ExecStart      = java -jar  /mnt/github/tjf/freshcrafts/wiz_postgres/target/wiz_postgres-0.0.1-SNAPSHOT.jar

[Install]
WantedBy = multi-user.target
